# root/etc/services.d/mpd/run
#!/usr/bin/with-contenv bash

exec mpd --stdout --no-daemon
