# root/etc/services.d/ympd/run
#!/usr/bin/with-contenv bash

exec ympd -h localhost -p 6600 -w 64801
